# supervisord.conf
; =========================================================================
; Dreampick — supervisord configuration for the Docker container.
; Boots the FastAPI backend and the React production build in parallel.
; =========================================================================

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; ------------------------- FastAPI backend -------------------------
[program:backend]
command=python backend/server.py
directory=/app
autostart=true
autorestart=true
startsecs=5
startretries=5
stopsignal=TERM
stopwaitsecs=15
environment=PYTHONUNBUFFERED="1",BACKEND_HOST="0.0.0.0",BACKEND_PORT="8000"
stdout_logfile=/var/log/supervisor/backend.out.log
stderr_logfile=/var/log/supervisor/backend.err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

; ------------------------- React static site (serve) -------------------------
[program:frontend]
command=serve -s frontend/build -l 3000
directory=/app
autostart=true
autorestart=true
startsecs=3
startretries=5
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/var/log/supervisor/frontend.out.log
stderr_logfile=/var/log/supervisor/frontend.err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
